# supervisord.conf
[supervisord]
nodaemon=false
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:cedarstar]
command=python main.py
directory=/opt/cedarstar
autostart=true
autorestart=true
startsecs=10
startretries=3
stderr_logfile=/var/log/supervisor/cedarstar.err.log
stdout_logfile=/var/log/supervisor/cedarstar.out.log
environment=PATH="/opt/cedarstar/venv/bin:%(ENV_PATH)s"
user=root

[program:repair_bot]
# 告诉 supervisor 咱们维修站代码所在的目录
directory=/opt/cedarstar
# 启动维修站的命令
command=python3 repair_bot.py
# 跟随服务器开机自启动
autostart=true
# 崩溃了自动拉起
autorestart=true
# 把报错和正常输出的日志存下来，方便宝宝随时查阅
stderr_logfile=/var/log/repair_bot.err.log
stdout_logfile=/var/log/repair_bot.out.log
# 如果你的 VPS 是用 root 登录的，就写 root
user=root
# 与 cedarstar 一致，优先使用项目 venv（python-telegram-bot 等依赖装在此处）
environment=PATH="/opt/cedarstar/venv/bin:%(ENV_PATH)s"
